# templates/webhook.service.erb
[Unit]
Description=Compose Webhook service

[Service]
User=deploy
Environment="WEBHOOK_JWT_SECRET=<%= @secret %>"
Environment="CONFIG_PATH=<%= @config_path %>"
ExecStart=/bin/bash -c "source ~/.rvm/scripts/rvm; compose-hook"
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target
